# supervisor-obs.conf
[program:obs]
priority=99
command=/usr/bin/obs --startstreaming
user=root
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
environment=DISPLAY=":1",HOME="/root",USER="root"