# supervisor-obs.conf
[program:obs]
priority=99
command=/usr/bin/obs --startstreaming
user=root
environment=DISPLAY=":1",HOME="/root",USER="root"